# go/foodscaner.service.example
[Unit]
Description=foodscaner
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/foodscaner/foodscaner
Restart=always
RestartSec=5

Environment=PORT=8080
Environment=ANTHROPIC_API_KEY=sk-ant-...
Environment=DB_PATH=/data/foodscaner.db
Environment=AUTH_URL=https://auth-center.sh-development.ru
Environment=AUTH_INTERNAL=https://auth-center.sh-development.ru
Environment=APP_URL=https://your-app-domain.com
Environment=APP_TOKEN=your-app-token-here
Environment=SECRET_KEY=your-secret-key-here

[Install]
WantedBy=multi-user.target
